# supervisor/army2-ci.conf
[program:army2-server]
directory=/root/army2
command=/usr/bin/java -jar army2.jar
autostart=true
autorestart=true
stdout_logfile=/root/log/army2-service.log
stderr_logfile=/root/log/army2_err.log
